# etc/supervisor.d/nginx.ini
[program:nginx]
command=/usr/sbin/nginx
stdout_logfile=/dev/null
stderr_logfile=/dev/null
stdout_events_enabled = true
stderr_events_enabled = true
