# etc/supervisor.d/uwsgi.ini
[program:uwsgi]
command=/usr/sbin/uwsgi --ini /etc/uwsgi/uwsgi.ini --master --die-on-term
stdout_logfile=/dev/null
stderr_logfile=/dev/null
stdout_events_enabled = true
stderr_events_enabled = true
